# crontabs/exports-crontab
SHELL=/bin/bash
WORKING_DIR=/home/blueskai/bluesky-ai-analysis
DATABASE_PATH=/mnt/ingestion/database/bluesky_posts.db
EXPORT_DIR=/mnt/ingestion/exports
VIRTUAL_ENV=.venv
VIRTUAL_ENV_HF=.venv-torch
LOG_DIR=/mnt/ingestion/logs

# Export embeddings to parquet every day at 11am and 11pm
0 11,23 * * * cd $WORKING_DIR && source $VIRTUAL_ENV/bin/activate && python scripts/export_embeddings.py --db-path $DATABASE_PATH --output-dir $EXPORT_DIR >> $LOG_DIR/export.log 2>&1
